# server/hirukara.supervisor.ini
[program:dosca]
command=/home/agile/p5-Hirukara-Lite/server/hirukara.runner.sh
directory=/home/agile/p5-Hirukara-Lite/
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/agile/p5-Hirukara-Lite/log/stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
